# packaging/systemd/deploybot-agent.service
[Unit]
Description=DeployBot Agent
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
EnvironmentFile=-/etc/deploybot-agent.env
ExecStart=/usr/local/bin/deploybot-agent --data-dir=/var/lib/deploybot --work-dir=/var/lib/deploybot/work
Restart=always
RestartSec=5
User=deploybot
Group=deploybot
SupplementaryGroups=docker
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
CapabilityBoundingSet=
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
